# crontab
* * * * * curl http://localhost:5555/ping > /dev/null 2>&1
